# include/nvidia/jetson/rc.local
#!/bin/sh
sh /root/enable_CAN.sh &
jetson_clocks
exit 0
